# debian/postinst
#!/bin/sh
# postinst script for argonone
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
	  #echo "Building New Virtual Environment in /opt/argonone..."
      #python3 -m venv --clear /opt/argonone
      #/opt/argonone/bin/python3 -m pip install setuptools wheel 
	  #/opt/argonone/bin/python3 -m pip install pysmbus RPi.GPIO 

      echo "Enabling i2C..."
      BOOT_CFG="/boot/config.txt"
      if [ -e "${BOOT_CFG}" ]; then
        touch "${BOOT_CFG}"
      fi
      if ! grep "dtparam=i2c_arm=on" "${BOOT_CFG}" >/dev/null 2>&1 ; then
        echo "dtparam=i2c=on" >> "${BOOT_CFG}"
      fi
      if ! grep "dtparam=i2c-1=on" "${BOOT_CFG}" >/dev/null 2>&1 ; then
        echo "dtparam=i2c-1=on" >> "${BOOT_CFG}"
      fi

      MODULE_CFG="/etc/modules-load.d/rpi2.conf"
      if [ -e "${MODULE_CFG}" ]; then
        touch "${MODULE_CFG}"
        chmod 644 "${MODULE_CFG}"
      fi
      if ! grep "i2c-dev" "${MODULE_CFG}" >/dev/null 2>&1 ; then
        echo "i2c-dev" >> "${MODULE_CFG}"
      fi
      if ! grep "i2c-bcm2708" "${MODULE_CFG}" >/dev/null 2>&1 ; then
        echo "i2c-bcm2708" >> "${MODULE_CFG}"
      fi
	  echo "Please reboot before first use."
      echo "Configuration file: /etc/argononed.conf"
      echo "Enable power button and fan control service with: sudo systemctl enable argononed.service"
      echo "If you want to: enable_uart=1"
      echo "Put it in /boot/config.txt and remove kgdboc=ttyAMA0,115200 in /boot/cmdline"
      echo "Thanks to thegeek1977 to help patching this"
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
